# .nerv/systemd/n3rv-hub.service
[Unit]
Description=N3RV A2A Hub for n3rv
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ultravietnamita/Projects/nerv
ExecStart=/home/ultravietnamita/Projects/nerv/.venv/bin/n3rv hub start
Restart=on-failure
RestartSec=5
StandardOutput=append:/home/ultravietnamita/Projects/nerv/.nerv/logs/hub.log
StandardError=append:/home/ultravietnamita/Projects/nerv/.nerv/logs/hub.log

[Install]
WantedBy=default.target
